# office-ai-detector-out.service
[Unit]
Description=Office AI Surveillance Detector OUT Service
Documentation=https://github.com/yourorg/office-ai-surveillance
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trinetra-jetson
Group=trinetra-jetson
WorkingDirectory=/home/trinetra-jetson/office-ai-surveillance
EnvironmentFile=-/etc/default/office-ai-detector-out

# Start OUT detector with auto-restart
ExecStart=/home/trinetra-jetson/office-ai-surveillance/detector-start-out.sh

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=office-ai-detector-out

# Resource limits
MemoryLimit=2G

# Security (relaxed for camera access)
NoNewPrivileges=false
ProtectHome=false
ReadWritePaths=/home/trinetra-jetson/office-ai-surveillance/python-scripts/crossings

[Install]
WantedBy=multi-user.target
